# server_on_ubuntu/server/MySQL/crontab
@reboot /home/odroid/Desktop/server/MySQL/mqtt_mysql_handler.sh
*/15 * * * * /home/odroid/Desktop/server/MySQL/mqtt_mysql_handler.sh

