# server_on_ubuntu/MySQL/crontab
@reboot /home/odroid/Desktop/server/mqtt_mysql_handler.sh
*/15 * * * * /home/odroid/Desktop/server/mqtt_mysql_handler.sh

